# config/superisor.conf
[program:app-uwsgi]
command = uwsgi restart

[program:nginx-app]
command = service nginx restart
